# usr/lib/opkg/info/ti-tisdk-loadmodule.postrm
#!/bin/sh
if test "x$D" != "x"; then
	OPT="-r $D"
else
	OPT=""
fi
if test "$1" = "remove" -o "$1" = "purge"; then
	if ! test -e "/etc/init.d/loadmodule-rc"; then
		update-rc.d $OPT loadmodule-rc remove
	fi
fi
